# systemd/falcon.service
[Unit]
Description=Falcon Auto-Commit Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=elixpo_workstation
WorkingDirectory=/home/elixpo_workstation/falcon
ExecStart=/home/elixpo_workstation/falcon/venv/bin/python3 /home/elixpo_workstation/falcon/falcon.py
StandardOutput=append:/home/elixpo_workstation/falcon/falcon.log
StandardError=append:/home/elixpo_workstation/falcon/falcon.log
Environment=HOME=/home/elixpo_workstation
Environment=GIT_AUTHOR_NAME=Circuit-Overtime
Environment=GIT_AUTHOR_EMAIL=ayushbhatt633@gmail.com
Environment=GIT_COMMITTER_NAME=Circuit-Overtime
Environment=GIT_COMMITTER_EMAIL=ayushbhatt633@gmail.com

[Install]
WantedBy=multi-user.target
